# z3/systemd.units/test-backup.target
[Unit]
Description=target for backup.sh
Documentation=

[Install]
WantedBy=multi-user.target
